# breaktimer.service
[Unit]
Description=Break timer

[Service]
Type=simple
StandardOutput=journal
ExecStart=<script>

[Install]
WantedBy=default.target
